# environments/dev/modules/private/files/supervisord-xrdp.conf
[program:xrdp]
priority = 15
command = /usr/sbin/xrdp --nodaemon
autostart = true
startsecs = 0
startretries = 0
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/xrdp
stdout_events_enabled = true
